# src/supplierservice/crontab
# START CRON JOB
*/5 * * * * /bin/bash -c "python3 /app/supplier_payments_notifier.py" > /proc/1/fd/1 2>/proc/1/fd/2
# END CRON JOB